# deploy/systemd/splitbot.service
[Unit]
Description=LINE Split Bot (Flask + gunicorn)
After=network.target

[Service]
Type=simple
User=splitbot
Group=splitbot
WorkingDirectory=/opt/splitbot
EnvironmentFile=/opt/splitbot/.env
# If your project path contains spaces, wrap the executable in double quotes:
#   ExecStart="/path/with space/.venv/bin/gunicorn" --bind 127.0.0.1:5050 ...
ExecStart=/opt/splitbot/.venv/bin/gunicorn --bind 127.0.0.1:5050 --workers 2 --access-logfile - --error-logfile - app:app
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
